# manifest/panel/etc/crontabs/nginx
* * * * * php /var/www/html/artisan p:schedule:process >> /dev/null 2>&1
