# manifest/etc/supervisor.d/runtime-files/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm7 --nodaemonize -c /etc/php7
autostart=true
autorestart=true
priority=5
